# src/tests/docker/issue94/faup-test.cron
*/1 * * * * /root/runfaup.sh > /tmp/faup.out
